# ops/systemd/hermes-gateway.service
[Unit]
Description=Hermes Gateway for htops
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/root/htops
Environment=HOME=/root
Environment=PATH=/root/.volta/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
Environment=HETANG_ROOT_DIR=/root/htops
EnvironmentFile=-/root/htops/.env.runtime
ExecStart=/usr/bin/env bash /root/htops/ops/hermes-gateway.sh
Restart=always
RestartSec=5
TimeoutStopSec=20

[Install]
WantedBy=multi-user.target
