# ops/systemd/htops-tools.service
[Unit]
Description=Hetang Hermes Tools Facade
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/htops
Environment=HETANG_ROOT_DIR=%h/htops
Environment=HETANG_RUNTIME_ENV_FILE=%h/htops/.env.runtime
ExecStart=/usr/bin/env bash %h/htops/ops/htops-tools.sh
Restart=always
RestartSec=5
StandardOutput=append:/tmp/htops-tools.systemd.log
StandardError=append:/tmp/htops-tools.systemd.log

[Install]
WantedBy=multi-user.target
